# infra/systemd/opsly-orchestrator-worker.service.example
# Copiar a /etc/systemd/system/opsly-orchestrator-worker.service y ajustar User/WorkingDirectory.
# Instalación: sudo cp infra/systemd/opsly-orchestrator-worker.service.example /etc/systemd/system/opsly-orchestrator-worker.service
#             sudo systemctl daemon-reload && sudo systemctl enable --now opsly-orchestrator-worker
#
# Requiere REDIS_URL en /home/opslyquantum/opsly/.env.local (o EnvironmentFile sin comillas raras).

[Unit]
Description=Opsly Orchestrator Worker (BullMQ)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=opslyquantum
Group=opslyquantum
WorkingDirectory=/home/opslyquantum/opsly
# Carga nvm + .env.local vía script
ExecStart=/home/opslyquantum/opsly/scripts/run-worker-with-nvm.sh
Restart=on-failure
RestartSec=20

[Install]
WantedBy=multi-user.target
