# infra/systemd/opsly-worker.service
# Instalación (en el worker Ubuntu, con sudo):
#   sudo cp infra/systemd/opsly-worker.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now opsly-worker
#
# Requiere: ~/opsly/.env.local con REDIS_URL; scripts/run-worker-with-nvm.sh ejecutable.
# run-worker-with-nvm.sh ejecuta git-sync-repo.sh antes de Node (pull --ff-only). Emergencia: OPSLY_SKIP_GIT_PULL=1 en Environment= del unit.

[Unit]
Description=Opsly Orchestrator Worker (BullMQ)
Documentation=https://github.com/cloudsysops/opsly
After=network-online.target
# Opcional si instalaste opsly-ollama.service: After=opsly-ollama.service
Wants=network-online.target

[Service]
Type=simple
User=opslyquantum
Group=opslyquantum
WorkingDirectory=/home/opslyquantum/opsly
ExecStart=/home/opslyquantum/opsly/scripts/run-worker-with-nvm.sh
Restart=always
RestartSec=15
StartLimitIntervalSec=300
StartLimitBurst=5

StandardOutput=append:/home/opslyquantum/opsly/logs/worker.log
StandardError=append:/home/opslyquantum/opsly/logs/worker-error.log

Environment=NODE_ENV=production

TimeoutStartSec=120
TimeoutStopSec=40

[Install]
WantedBy=multi-user.target
