# infra/systemd/opsly-mac2011-cleanup.timer
[Unit]
Description=Opsly Mac 2011 cleanup timer (Linux)

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
