# infra/systemd/opsly-disk-maintain-fanout.user.service
[Unit]
Description=Opsly disk threshold maintain and fleet fan-out (user session; worker)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Ajustar %h/opsly a la ruta real del clon en el worker (p. ej. %h/proyectos/intcloudsysops).
WorkingDirectory=%h/opsly
EnvironmentFile=-%h/.config/opsly/disk-maintain.env
ExecStart=/bin/bash %h/opsly/scripts/opsly-disk-maintain-fanout.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
