# infra/launchd/com.opsly.mac2011.cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.opsly.mac2011.cleanup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__OPSLY_ROOT__/scripts/mac2011-cleanup-robust.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>OPSLY_LOG_DIR</key>
    <string>__OPSLY_ROOT__/logs</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__OPSLY_ROOT__/logs/launchd-mac2011-cleanup.out</string>
  <key>StandardErrorPath</key>
  <string>__OPSLY_ROOT__/logs/launchd-mac2011-cleanup.err</string>
</dict>
</plist>
